# root/etc/s6-overlay/services.d/web/run
#!/usr/bin/with-contenv sh
cd /
exec /app/web
